# deploy/detect-backend.service
[Unit]
Description=Detect App Backend (FastAPI + YOLOv8)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/detect-app/backend
Environment=PATH=/var/www/detect-app/backend/venv/bin
Environment=PYTHONPATH=/var/www/detect-app/backend
ExecStart=/var/www/detect-app/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
